# services/autoscript.service
[Unit]
Description=Auto script
 
[Service]
ExecStart=/home/sms/scripts/autoscript.sh
 
[Install]
WantedBy=default.target

